# curator4-881992795.conf
[program:curator4-881992795]
command=python3 /home/curator4/curator4-881992795/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-881992795/test.err.log
stdout_logfile=/home/curator4/curator4-881992795/test.out.log